# deployment/fisk-client/DEBIAN/postinst
cd /var/fisk/daemon
sudo -H npm install --unsafe-perm @andersbakken/fisk@latest
ln -s /var/fisk/daemon/node_modules/@andersbakken/fisk/monitor/fisk-monitor.js /usr/local/bin/fisk-monitor
mv /usr/local/bin/fiskc.linux /usr/local/bin/fiskc
rm /usr/local/bin/fiskc.osx
if [ -x "`which systemctl`" ] && ! uname -r | grep -i microsoft; then
   systemctl daemon-reload
   systemctl enable fisk-daemon.socket
   systemctl enable fisk-daemon.service
   # Start (or restart) the socket first so the listening fd survives
   # daemon upgrades. Existing clients (including bind-mounted docker
   # containers) keep working across restarts of fisk-daemon.service.
   systemctl restart fisk-daemon.socket
   systemctl restart fisk-daemon.service
fi

mkdir -p /var/fisk/client/compilers
chmod 777 /var/fisk/client/compilers

if [ ! -x "/usr/local/bin/fisk-monitor" ]; then
    MONITOR_DIR=/usr/local/fisk/monitor
    LINK_SRC="$MONITOR_DIR/node_modules/@andersbakken/fisk/monitor/fisk-monitor.js"
    mkdir $MONITOR_DIR
    cd $MONITOR_DIR
    npm init -y
    npm install --save @andersbakken/fisk@latest
    ln -s "$LINK_SRC" /usr/local/bin/fisk-monitor
fi
