# deployment/fisk-scheduler/DEBIAN/postinst
#!/bin/bash

cd /var/fisk/scheduler
sudo -H npm install --unsafe-perm
sudo systemctl enable fisk-scheduler.service
sudo service fisk-scheduler start

