# deployment/fisk-scheduler/DEBIAN/postinst
#!/bin/bash

mkdir -p /var/fisk/scheduler/data
npm install --unsafe-perm -g @andersbakken/fisk
sudo systemctl enable fisk-scheduler.service
sudo service fisk-scheduler start

