# deployment/fisk-scheduler/DEBIAN/prerm
#!/bin/bash

rm -rf /var/fisk/scheduler/data
sudo service fisk-scheduler stop
npm  --unsafe-perm uninstall -g @andersbakken/fisk
sudo systemctl disable fisk-scheduler.service
